# debian/DEBIAN/prerm
#!/bin/bash

systemctl stop utsms-daemon || true
systemctl disable utsms-daemon || true

exit 0
